# test/supervisord.conf
[supervisord]
nodaemon=true

[program:player]
command=python player.py
directory=/app
autorestart=true
startretries=10

[program:master]
command=python master.py
directory=/app
autorestart=true
startretries=10
